# crypto-trading.service
[Unit]
Description=crypto trading daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/.openclaw/workspace/projects/crypto-trading
ExecStart=/root/.openclaw/workspace/projects/crypto-trading/run_daemon.sh
Restart=always
RestartSec=5
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
