# deploy/systemd/crypto-trading-research-pipeline.service
[Unit]
Description=crypto-trading research pipeline
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/.openclaw/workspace/projects/crypto-trading
ExecStart=/root/.openclaw/workspace/projects/crypto-trading/scripts/pipeline/run_research_pipeline.sh
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
